# roles/audio/templates/audio/playerctld.service.j2
[Unit]
Description=Playerctl daemon
PartOf=default.target
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'playerctld'
Restart=on-failure

[Install]
WantedBy=default.target
