# templates/profile.d/zz-cleanup.sh.j2
# Cleanup
unset -f prepend_envvar
unset -f import_env
